# systemd/user/waybar.service
[Unit]
Description=Start Waybar

[Service]
Type=oneshot
ExecStart=/sbin/waybar

[Install]
WantedBy=default.target
